# .config/crontab
* * * * * /home/user/.local/bin/battery_checker
